# application/docker/supervisor.conf
[supervisord]
nodaemon=true


[program:executablejar]
directory = /
command = java -jar app.jar
autostart = true
autorestart = false